# contrib/systemd/scalaris.service.in
[Unit]
Description=Scalaris node
Requires=remote-fs.target
After=remote-fs.target

[Service]
EnvironmentFile=@sysconfdir@/conf.d/scalaris
ExecStart=@bindir@/scalarisctl -n ${SCALARIS_NODE} -e "-noshell" -l @localstatedir@/log/scalaris $SCALARIS_ADDITIONAL_PARAMETERS start
ExecStop=@bindir@/scalarisctl -n ${SCALARIS_NODE} -e "-noshell" -l @localstatedir@/log/scalaris $SCALARIS_ADDITIONAL_PARAMETERS gstop

WorkingDirectory=@localstatedir@/log/scalaris

User=scalaris
Group=scalaris
