# contrib/systemd/scalaris.service.in
[Unit]
Description=Scalaris node
Requires=remote-fs.target
After=remote-fs.target

[Service]
EnvironmentFile=@sysconfdir@/conf.d/scalaris
ExecStart=@bindir@/scalarisctl -n ${SCALARIS_NODE} -s -e "-noshell" start
ExecStop=@bindir@/scalarisctl -n ${SCALARIS_NODE} -s -e "-noshell" gstop

WorkingDirectory=@localstatedir@/log/scalaris

User=scalaris
Group=scalaris
