# tags/0.5.0/contrib/packages/main/debian.scalaris.postrm
#!/bin/sh
set -e

if [ "$1" = "purge" ]; then
  rm -f /etc/scalaris/scalaris.cfg
  echo "Removing Scalaris RC-links"
  update-rc.d -f scalaris remove > /dev/null
fi

#DEBHELPER#

exit 0
